# systemd/refindsnapper.service
[Unit]
Description=Regenerate rEFInd snapper
StartLimitBurst=0
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/usr/bin/refindsnapper
