# systemd/refindsnapper.path
[Unit]
Description=Monitors for new snapshots

[Path]
PathChanged=/.snapshots
PathChanged=/efi/loader/entries

[Install]
WantedBy=multi-user.target
